# initramfs_root/etc/inittab
#::sysinit:/etc/load_root_from_flash.sh
::askfirst:-/bin/sh
#::respawn:/sbin/getty -L 115200 tty0 vt100
::restart:/sbin/init
::shutdown:/bin/umount -a -r
